# Docker_files/Ubuntu2/root/etc/services.d/mongod/run
#!/usr/bin/with-contenv bash

echo "Starting MongoDB"

exec /usr/bin/mongod --config /etc/mongodb.conf run